# scripts/luvora-cron.service
# Systemd Service for Luvora Spark Delivery
# Installation:
#   sudo cp luvora-cron.service /etc/systemd/system/
#   sudo cp luvora-cron.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable luvora-cron.timer
#   sudo systemctl start luvora-cron.timer

[Unit]
Description=Luvora Spark Delivery Service
After=network.target

[Service]
Type=oneshot
User=your-user
WorkingDirectory=/path/to/Luvora
Environment="APP_URL=http://localhost:3000"
Environment="CRON_SECRET=Luvora"
ExecStart=/path/to/Luvora/scripts/cron-deliver.sh

# Logging
StandardOutput=append:/var/log/luvora/cron-deliver.log
StandardError=append:/var/log/luvora/cron-deliver-error.log

# Security
PrivateTmp=true
NoNewPrivileges=true

# Restart on failure
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
